# launchd/com.example.sirius-fetch.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  全ソース一括取得を1日3回（7:00 / 12:00 / 19:00）。スリープ中は起床後に実行される。
  セットアップ: __INSTALL_DIR__ と __HOME__ を自分の環境の絶対パスに置換し、
  ~/Library/LaunchAgents/ に置いて `launchctl load` する。
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.sirius-fetch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>__INSTALL_DIR__/fetch-all.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <key>WorkingDirectory</key>
  <string>__INSTALL_DIR__</string>

  <key>StandardOutPath</key>
  <string>__HOME__/.local/share/sirius-fetch/fetch.out.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.local/share/sirius-fetch/fetch.err.log</string>
</dict>
</plist>
